# ubuntu/apps/keyboard-config/keyd/systemd/keyd.timer
[Unit]
Description=key remapping daemon start timer

[Timer]
OnBootSec=5

[Install]
WantedBy=default.target
